# overlay/var/cron/crontabs/root
*/5 * * * * /maasland_app/scripts/flexess_cron.php
