# overlay/var/cron/crontabs/root
*/3 * * * * /maasland_app/scripts/flexess_cron.sh
0 4 * * * /maasland_app/scripts/maintenance.sh
